# sfa/cron.d/sfa.cron
0 * * * *  /usr/bin/sfa-import-plc.py > /dev/null 2>&1
0 0 * * * /usr/bin/sfa-clean-peer-records.py > /dev/null 2>&1  
